# jupyterhub.conf
# this file should be in /etc/init/jupyterhub.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn
expect fork
setuid jupyterhub
setgid shadow
exec /software/jupyterhub/launchjupyterhub.sh
